# rootfs/etc/services.d/beamer-client/run
#!/bin/bash

echo "Starting USB Beamer Client service..."
python3 /beamer_client/main.py 